# systemd/cert-server.service.in
[Unit]
Description=Start cert server for cert-svc managing root ssl certs

[Service]
ExecStart=@TZ_SYS_BIN@/cert-server
User=security_fw
Group=security_fw
Sockets=cert-server.socket

[Install]
WantedBy=multi-user.target
